# .systemd/hauski-reviewd.timer
[Unit]
Description=hausKI review daemon timer

[Timer]
OnBootSec=2min
OnUnitActiveSec=5min
Unit=hauski-reviewd.service
AccuracySec=30s

[Install]
WantedBy=timers.target
